# config/systemd/sidekiq.service.erb
[Unit]
Description = <%= fetch(:application) %> Sidekiq Background Worker
After = multi-user.target

[Service]
Type = simple
Environment = RAILS_ENV=<%= fetch(:rails_env) %>
Environment = PWD=<%= current_path %>
EnvironmentFile=/etc/environment
WorkingDirectory = <%= current_path %>

User = <%= fetch(:user) %>
Group = <%= fetch(:group) || fetch(:user) %>
KillMode = process
UMask = 0002
RestartSec = 1
Restart = always

[Install]
WantedBy = multi-user.target
